# debian/postinst
#!/bin/bash
systemctl enable switch_camera
systemctl start switch_camera
